# com.auctionwatcher.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.auctionwatcher.api</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/claw/auction-watcher/api/server.py</string>
        <string>--port</string>
        <string>7474</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/claw/auction-watcher</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/claw/auction-watcher/logs/api_server_launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/claw/auction-watcher/logs/api_server_launchd.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/claw</string>
    </dict>
</dict>
</plist>
